# crontab
# build image
0 1 * * *   python3 build_script.py build 10

# run container
0 5 * * *   python3 build_script.py run 10